# examples/play-example/src/codedeploy/content/etc/init/play-example.conf
start on started elastic-network-interfaces
respawn
script
	/usr/share/play-example/bin/play-example 2>&1 >> /var/log/play-example.log
end script
